# root/etc/s6-overlay/s6-rc.d/svc-fetchbox/run
#!/usr/bin/with-contenv bash
exec s6-setuidgid abc /usr/local/bin/fetchbox --daemon --config /config/fetchbox.yml
